# deploy/github-app-token-refresher.service
[Unit]
Description=Refresh TinyAssets GitHub App installation token
Documentation=https://github.com/Jonnyton/TinyAssets/blob/main/docs/ops/bot-identity-setup.md
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tinyassets/github-app-token-refresher.env

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/tinyassets
EnvironmentFile=/etc/tinyassets/github-app-token-refresher.env
ExecStart=/usr/bin/python3 /opt/tinyassets/scripts/github-app-token-refresher.py
TimeoutStartSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-app-token-refresher
