# deploy/tinyassets-backup.service
# TinyAssets daemon state backup — oneshot snapshot + upload.
#
# Self-host migration Row J per
# docs/exec-plans/active/2026-04-20-selfhost-uptime-migration.md.
#
# Triggered by tinyassets-backup.timer nightly at 03:00 UTC. Invokes
# deploy/backup.sh: brain-tier + full-tier archives of the tinyassets-data
# named volume, rclone-uploaded to BACKUP_DEST and shipped to GitHub
# release assets when GH_TOKEN is set. Retention pruned in the same run.
#
# Runs as root because:
#   1. /var/lib/docker/volumes/ is only root-readable by default.
#   2. `docker volume inspect` needs docker socket access — simpler
#      to grant via root than extend the tinyassets user's docker-group
#      membership across every ops primitive.
#
# Backup destination and retention settings live in /etc/tinyassets/env (the
# same env file the daemon uses). systemd sources it via EnvironmentFile;
# destination credentials live in root's rclone configuration.

[Unit]
Description=TinyAssets daemon state backup (nightly two-tier snapshot to BACKUP_DEST + GH releases)
Documentation=https://github.com/Jonnyton/TinyAssets/blob/main/deploy/HETZNER-DEPLOY.md
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/tinyassets/env
ExecStart=/bin/bash /opt/tinyassets-host-uptime/current/deploy/backup.sh
TimeoutStartSec=30min
# Don't kill mid-upload — backup.sh handles its own timeouts.

# Isolate temporary archives from other services. Root's persistent rclone
# configuration remains outside PrivateTmp at /root/.config/rclone/.
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyassets-backup
