# deploy/tinyassets-disk-watch.service
# TinyAssets disk pressure watchdog — hourly disk check + auto-prune.
#
# Runs three steps per tick (BUG-023 Phases 2+3+4):
#   1. disk_watch.py   — open GH issue if /var/lib/docker is past DISK_WARN_PCT.
#   2. rotate_run_transcripts — gzip + archive old run artifacts.
#   3. disk_autoprune.py — docker image prune -af if root ≥ DISK_AUTOPRUNE_PCT.
#
# The ExecStart commands run sequentially. Status 1 is accepted because
# disk_watch.py intentionally returns it under pressure; other statuses fail
# the unit. Rotation and auto-prune handle their documented optional paths.
#
# Environment variables sourced from /etc/tinyassets/env:
#   GITHUB_TOKEN       — fine-grained PAT with issues:write
#   GITHUB_REPOSITORY  — owner/repo  (default: Jonnyton/TinyAssets)
#   DISK_WARN_PCT      — disk_watch GH-alert threshold (default 80)
#   DISK_AUTOPRUNE_PCT — prune trigger threshold (default 85)
#   DISK_WATCH_PATH    — path to check, default /var/lib/docker for
#                        disk_watch, / for disk_autoprune
#   TINYASSETS_RUN_TRANSCRIPT_RETENTION_DAYS — rotation retention (default 30)
#
# Install:
#   sudo cp deploy/tinyassets-disk-watch.service /etc/systemd/system/
#   sudo cp deploy/tinyassets-disk-watch.timer    /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tinyassets-disk-watch.timer

[Unit]
Description=TinyAssets disk watchdog + rotation + auto-prune (hourly)
Documentation=https://github.com/Jonnyton/TinyAssets/blob/main/docs/ops/log-aggregation-runbook.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
SuccessExitStatus=1
EnvironmentFile=/etc/tinyassets/env
# Step 1 — legacy disk_watch GH-issue alarm. Its intentional pressure status
# 1 is accepted above so rotation and auto-prune can still run.
ExecStart=/usr/bin/python3 /opt/tinyassets/scripts/disk_watch.py
# Step 2 — rotate run transcripts older than retention into archived/.
ExecStart=/usr/bin/python3 -m scripts.rotate_run_transcripts
# Step 3 matches p0-outage disk_full cleanup: Docker system/builder
# prune plus journal vacuum, without Docker volumes.
# Step 3 — auto-prune docker images if root volume is pressured.
ExecStart=/usr/bin/python3 /opt/tinyassets/scripts/disk_autoprune.py
# Accepted ExecStart statuses continue in order; an unexpected status stops
# later commands. Type=oneshot treats the sequence as one timer activation.
TimeoutStartSec=180s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyassets-disk-watch
